# container_scripts/supervisord-hiawatha.conf
[program:hiawatha]
command=/start.sh
numprocs=1
autostart=true
autorestart=true

